# scripts/webcface-launcher.service.in
[Unit]
Description=WebCFace Launcher

[Service]
Type=exec
Restart=always
ExecStart=@prefix@/@bindir@/webcface-launcher @WEBCFACE_TOOLS_LAUNCHER_CONFIG_PATH@

[Install]
WantedBy=multi-user.target
